# nixpacks/config/workers/laravel-queue-worker.conf
; Supervisor configuration for Laravel Queue Worker
; Processes background jobs from the queue

[program:laravel-queue-worker]
process_name = %(program_name)s_%(process_num)02d
directory = /app
command = php artisan queue:work --sleep=3 --tries=3 --max-time=3600
user = www-data
numprocs = 1
autostart = false
autorestart = true
stopasgroup = true
killasgroup = true
startsecs = 0
stopwaitsecs = 3600
priority = 30
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
