# nixpacks/config/workers/nginx.conf
; Supervisor configuration for Nginx process
; Manages Nginx web server lifecycle

[program:nginx]
process_name = %(program_name)s
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
user = www-data
autostart = true
autorestart = true
startsecs = 0
stopasgroup = true
killasgroup = true
priority = 10
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
