# nixpacks/config/workers/php-fpm.conf
; Supervisor configuration for PHP-FPM process
; Manages PHP-FPM lifecycle for handling PHP requests

[program:php-fpm]
process_name = %(program_name)s
command = php-fpm -y /etc/php-fpm/php-fpm.conf -F
user = root
autostart = true
autorestart = true
startsecs = 0
stopasgroup = true
killasgroup = true
priority = 20
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
